# extra/ranger-vmware-executor.service
[Unit]
Description=ranger-vmware-executor
After=network.target

[Service]
ExecStart=/var/opt/ranger/bin/ranger-vmware-executor /etc/opt/ranger/ranger-vmware-executor/config.yml
Type=simple
User=root
Group=root
Restart=on-failure
RestartSec=10
StartLimitInterval=10m
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
